# scripts/strapi.service
[Unit]
Description=Strapi CMS
After=network.target

[Service]
Type=simple
User=nem
Group=nem
WorkingDirectory=/home/nem/code/NewEnglandMelee/strapi
# Load environment variables from .env file
EnvironmentFile=/home/nem/code/NewEnglandMelee/strapi/.env
ExecStart=/usr/bin/node /home/nem/code/NewEnglandMelee/strapi/dist/src/index.js
# Alternative: ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
MemoryMax=1G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
